# luci-app-taoistfuchen/root/etc/init.d/fakehttp
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

PROG="/usr/bin/fakehttp"
COMMON="/usr/share/taoistfuchen/service-common.sh"
SERVICE_NAME=fakehttp
DEFAULT_BOOT_DELAY=60
EXTRA_COMMANDS="link_down link_up boot_expired"
EXTRA_HELP="
	link_down DEVICE     Stop for a matching link-down event
	link_up DEVICE       Resume for a matching link-up event
	boot_expired TOKEN DEADLINE  Complete a supervised boot wait"

. "$COMMON"

fh_error() {
	logger -t fakehttp -- "$*"
}

fh_validate_interface() {
	local iface="$1"

	[ -n "$iface" ] || return 0
	if ! tf_valid_interface "$iface" || ! tf_interface_exists "$iface"; then
		fh_error "invalid or unavailable network interface: $iface"
		FH_VALID=0
		return 0
	fi
	case " $FH_INTERFACES " in *" $iface "*) return 0 ;; esac
	FH_INTERFACES="${FH_INTERFACES:+$FH_INTERFACES }$iface"
	FH_INTERFACE_COUNT=$((FH_INTERFACE_COUNT + 1))
}

fh_validate_payload() {
	local section="$1" enabled type host file

	config_get_bool enabled "$section" enabled 1
	[ "$enabled" = "1" ] || return 0
	config_get type "$section" type ''
	case "$type" in
		http|https)
			config_get host "$section" host ''
			if ! tf_valid_hostname "$host"; then
				fh_error "invalid $type payload hostname in section $section"
				FH_VALID=0
				return 0
			fi
			;;
		binary)
			config_get file "$section" file ''
			if ! tf_valid_fakehttp_payload_path "$file"; then
				fh_error "invalid, missing, or oversized binary payload in section $section"
				FH_VALID=0
				return 0
			fi
			;;
		*)
			fh_error "unsupported payload type in section $section"
			FH_VALID=0
			return 0
			;;
	esac
	FH_PAYLOAD_COUNT=$((FH_PAYLOAD_COUNT + 1))
}

fh_validate_config() {
	local interface_mode direction family repeat ttl hop_estimation dynamic_ttl_pct
	local queue_num sibling_queue fwmark fwmask log_connections

	FH_VALID=1
	FH_INTERFACE_COUNT=0
	FH_PAYLOAD_COUNT=0
	FH_INTERFACES=''
	config_get interface_mode main interface_mode 'selected'
	case "$interface_mode" in selected|all) ;; *) fh_error "invalid interface mode"; FH_VALID=0 ;; esac
	[ "$interface_mode" = "selected" ] && config_list_foreach main interface fh_validate_interface
	config_foreach fh_validate_payload payload

	if [ "$interface_mode" = "selected" ] && [ "$FH_INTERFACE_COUNT" -eq 0 ]; then
		fh_error "select at least one existing Internet-facing interface"
		FH_VALID=0
	fi
	if [ "$FH_PAYLOAD_COUNT" -eq 0 ]; then
		fh_error "configure at least one valid payload"
		FH_VALID=0
	fi

	config_get direction main direction 'outbound'
	case "$direction" in inbound|outbound|both) ;; *) fh_error "invalid direction"; FH_VALID=0 ;; esac
	config_get family main family 'dual'
	case "$family" in ipv4|ipv6|dual) ;; *) fh_error "invalid address family"; FH_VALID=0 ;; esac

	config_get repeat main repeat '2'
	tf_valid_uint_range "$repeat" 1 10 || { fh_error "repeat must be 1-10"; FH_VALID=0; }
	config_get ttl main ttl '3'
	tf_valid_uint_range "$ttl" 1 255 || { fh_error "TTL must be 1-255"; FH_VALID=0; }
	config_get_bool hop_estimation main hop_estimation 1
	config_get dynamic_ttl_pct main dynamic_ttl_pct '0'
	case "$dynamic_ttl_pct" in
		0) ;;
		*) tf_valid_uint_range "$dynamic_ttl_pct" 1 99 || { fh_error "dynamic TTL percentage must be 0-99"; FH_VALID=0; } ;;
	esac
	if [ "$hop_estimation" != "1" ] && [ "$dynamic_ttl_pct" != "0" ]; then
		fh_error "dynamic TTL requires hop estimation"
		FH_VALID=0
	fi
	config_get queue_num main queue_num '8970'
	tf_valid_uint_range "$queue_num" 1 65535 || { fh_error "NFQUEUE number must be 1-65535"; FH_VALID=0; }
	if [ "$(uci -q get fakesip.main.enabled)" = "1" ]; then
		sibling_queue="$(uci -q get fakesip.main.queue_num)"
		if tf_valid_uint_range "$queue_num" 1 65535 &&
		   tf_valid_uint_range "$sibling_queue" 1 65535 &&
		   [ "$queue_num" -eq "$sibling_queue" ]; then
			fh_error "NFQUEUE $queue_num is already assigned to enabled FakeSIP"
			FH_VALID=0
		fi
	fi
	config_get fwmark main fwmark '0x8000'
	config_get fwmask main fwmask '0x8000'
	tf_valid_u32 "$fwmark" || { fh_error "invalid fwmark"; FH_VALID=0; }
	tf_valid_u32 "$fwmask" || { fh_error "invalid fwmark mask"; FH_VALID=0; }
	if tf_valid_u32 "$fwmark" && tf_valid_u32 "$fwmask"; then
		tf_mark_fits_mask "$fwmark" "$fwmask" || { fh_error "fwmark must fit inside its mask"; FH_VALID=0; }
	fi
	config_get_bool log_connections main log_connections 0
	case "$log_connections:$hop_estimation" in
		0:0|0:1|1:0|1:1) ;;
		*) fh_error "invalid boolean option"; FH_VALID=0 ;;
	esac

	[ "$FH_VALID" = "1" ]
}

fh_append_interface() {
	local iface="$1"

	[ -n "$iface" ] || return 0
	case " $FH_APPENDED_INTERFACES " in *" $iface "*) return 0 ;; esac
	FH_APPENDED_INTERFACES="${FH_APPENDED_INTERFACES:+$FH_APPENDED_INTERFACES }$iface"
	procd_append_param command -i "$iface"
	procd_append_param netdev "$iface"
}

fh_append_payload() {
	local section="$1" enabled type host file

	config_get_bool enabled "$section" enabled 1
	[ "$enabled" = "1" ] || return 0
	config_get type "$section" type
	case "$type" in
		http)
			config_get host "$section" host
			procd_append_param command -h "$host"
			;;
		https)
			config_get host "$section" host
			procd_append_param command -e "$host"
			;;
		binary)
			config_get file "$section" file
			procd_append_param command -b "$file"
			;;
	esac
}

fh_collect_payload_reverse() {
	local section="$1"
	FH_REVERSED_PAYLOADS="$section${FH_REVERSED_PAYLOADS:+ $FH_REVERSED_PAYLOADS}"
}

fh_match_event_interface() {
	[ "$1" = "$FH_EVENT_DEVICE" ] && FH_EVENT_MATCH=1
}

fh_event_device_matches() {
	local device="$1" enabled interface_mode

	config_load fakehttp
	config_get_bool enabled main enabled 0
	[ "$enabled" = 1 ] || return 1
	config_get interface_mode main interface_mode selected
	case "$interface_mode" in
		all) return 0 ;;
		selected) ;;
		*) return 1 ;;
	esac
	FH_EVENT_DEVICE="$device"
	FH_EVENT_MATCH=0
	config_list_foreach main interface fh_match_event_interface
	[ "$FH_EVENT_MATCH" = 1 ]
}

fh_start_service_locked() {
	local enabled interface_mode direction family repeat ttl hop_estimation dynamic_ttl_pct
	local log_connections queue_num fwmark fwmask section
	local context boot_delay token deadline state

	config_load fakehttp
	config_get_bool enabled main enabled 0
	context="${TAOISTFUCHEN_SERVICE_CONTEXT:-}"
	if [ -z "$context" ]; then
		case "${action:-start}" in
			boot) context=boot ;;
			*) context=manual ;;
		esac
	fi

	case "$context" in
		ifdown-definition)
			return 0
			;;
		manual)
			if [ "$enabled" = 1 ]; then
				tf_boot_state_set "$SERVICE_NAME" manual manual 0 || return 1
			else
				tf_boot_state_set "$SERVICE_NAME" cancel disabled 0 || return 1
				( procd_kill "$SERVICE_NAME" ) >/dev/null 2>&1 || true
				tf_cleanup_nft_table "$SERVICE_NAME"
				return 0
			fi
			;;
		boot-expired)
			token="${TAOISTFUCHEN_BOOT_TOKEN:-}"
			deadline="${TAOISTFUCHEN_BOOT_DEADLINE:-}"
			tf_boot_token_valid "$token" || return 1
			tf_valid_uint_range "$deadline" 0 600 || return 1
			tf_boot_state_is "$SERVICE_NAME" wait "$token" "$deadline" || return 0
			if [ "$enabled" != 1 ]; then
				tf_boot_state_set "$SERVICE_NAME" cancel disabled 0 || return 1
				return 0
			fi
			tf_boot_state_set "$SERVICE_NAME" done "$token" "$deadline" || return 1
			tf_link_state_is_down "$SERVICE_NAME" && return 0
			;;
		boot|auto|ifup)
			config_get boot_delay main boot_delay "$DEFAULT_BOOT_DELAY"
			boot_delay="$(tf_boot_delay_value "$boot_delay" "$DEFAULT_BOOT_DELAY")" || return 1
			if [ "$enabled" != 1 ]; then
				tf_boot_state_set "$SERVICE_NAME" cancel disabled 0 || return 1
				( procd_kill "$SERVICE_NAME" ) >/dev/null 2>&1 || true
				tf_cleanup_nft_table "$SERVICE_NAME"
				return 0
			fi
			tf_boot_should_defer "$SERVICE_NAME" "$boot_delay" && return 0
			state="$(tf_boot_state_get "$SERVICE_NAME" 2>/dev/null || true)"
			case "$state" in
				manual\ *|done\ *) ;;
				*) tf_boot_state_set "$SERVICE_NAME" done automatic "$boot_delay" || return 1 ;;
			esac
			;;
		*)
			return 1
			;;
	esac

	( procd_kill "$SERVICE_NAME" ) >/dev/null 2>&1 || true
	tf_cleanup_nft_table "$SERVICE_NAME"
	[ -x "$PROG" ] || { fh_error "binary is missing or not executable: $PROG"; return 1; }
	fh_validate_config || return 1

	config_get interface_mode main interface_mode 'selected'
	config_get direction main direction 'outbound'
	config_get family main family 'dual'
	config_get repeat main repeat '2'
	config_get ttl main ttl '3'
	config_get_bool hop_estimation main hop_estimation 1
	config_get dynamic_ttl_pct main dynamic_ttl_pct '0'
	config_get_bool log_connections main log_connections 0
	config_get queue_num main queue_num '8970'
	config_get fwmark main fwmark '0x8000'
	config_get fwmask main fwmask '0x8000'

	procd_open_instance fakehttp
	procd_set_param command "$PROG"
	FH_APPENDED_INTERFACES=''
	if [ "$interface_mode" = "all" ]; then
		procd_append_param command -a
	else
		config_list_foreach main interface fh_append_interface
	fi
	case "$direction" in
		inbound) procd_append_param command -0 ;;
		outbound) procd_append_param command -1 ;;
		both) procd_append_param command -0 -1 ;;
	esac
	case "$family" in
		ipv4) procd_append_param command -4 ;;
		ipv6) procd_append_param command -6 ;;
		dual) procd_append_param command -4 -6 ;;
	esac
	# FakeHTTP 0.9.18 inserts each parsed payload immediately after the ring's
	# current node, so its effective rotation is the reverse of CLI order.
	# Reverse UCI sections here to make the observed rotation match the LuCI row
	# order while preserving duplicate rows as weights.
	FH_REVERSED_PAYLOADS=''
	config_foreach fh_collect_payload_reverse payload
	for section in $FH_REVERSED_PAYLOADS; do
		fh_append_payload "$section"
	done
	procd_append_param command -r "$repeat" -t "$ttl" -n "$queue_num" -m "$fwmark" -x "$fwmask"
	[ "$hop_estimation" = "1" ] || procd_append_param command -g
	[ "$dynamic_ttl_pct" = "0" ] || procd_append_param command -y "$dynamic_ttl_pct"
	[ "$log_connections" = "1" ] || procd_append_param command -s
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param respawn 3600 5 5
	procd_set_param term_timeout 5
	procd_close_instance
}

start_service() {
	local result=0

	tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
	fh_start_service_locked || result=$?
	tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
	return "$result"
}

link_down() {
	local device="${1:-}" result=0

	tf_valid_interface "$device" || { fh_error "invalid link-down device"; return 1; }
	procd_lock
	tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
	if ! fh_event_device_matches "$device"; then
		tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
		return 0
	fi
	tf_link_state_set "$SERVICE_NAME" down || result=$?
	if [ "$result" -eq 0 ]; then
		procd_kill "$SERVICE_NAME" >/dev/null 2>&1 || true
		tf_cleanup_nft_table "$SERVICE_NAME"
		TAOISTFUCHEN_SERVICE_CONTEXT=ifdown-definition start || result=$?
	fi
	tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
	return "$result"
}

link_up() {
	local device="${1:-}" result=0

	tf_valid_interface "$device" || {
		fh_error "invalid or unavailable link-up device: $device"
		return 1
	}
	procd_lock
	tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
	if ! tf_interface_exists "$device" || ! fh_event_device_matches "$device"; then
		tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
		return 0
	fi
	tf_link_state_set "$SERVICE_NAME" up || result=$?
	if [ "$result" -eq 0 ]; then
		procd_kill "$SERVICE_NAME" >/dev/null 2>&1 || true
		tf_cleanup_nft_table "$SERVICE_NAME"
		TAOISTFUCHEN_SERVICE_CONTEXT=ifup start || result=$?
	fi
	tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
	return "$result"
}

boot_expired() {
	local token="${1:-}" deadline="${2:-}" result=0

	tf_boot_token_valid "$token" || return 1
	tf_valid_uint_range "$deadline" 0 600 || return 1
	procd_lock
	tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
	if ! tf_boot_state_is "$SERVICE_NAME" wait "$token" "$deadline"; then
		tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
		return 0
	fi
	TAOISTFUCHEN_SERVICE_CONTEXT=boot-expired \
		TAOISTFUCHEN_BOOT_TOKEN="$token" \
		TAOISTFUCHEN_BOOT_DEADLINE="$deadline" \
		start || result=$?
	tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
	return "$result"
}

reload_service() {
	local result=0

	procd_lock
	tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
	procd_kill "$SERVICE_NAME" >/dev/null 2>&1 || true
	tf_cleanup_nft_table "$SERVICE_NAME"
	TAOISTFUCHEN_SERVICE_CONTEXT=manual start || result=$?
	tf_lifecycle_lock_release "$SERVICE_NAME" || return 1
	return "$result"
}

stop_service() {
	case "${action:-stop}" in
		stop|shutdown|restart)
			tf_lifecycle_lock_acquire "$SERVICE_NAME" || return 1
			tf_boot_state_set "$SERVICE_NAME" cancel explicit-stop 0 || {
				tf_lifecycle_lock_release "$SERVICE_NAME" || true
				return 1
			}
			tf_lifecycle_lock_release "$SERVICE_NAME"
			;;
	esac
}

service_triggers() {
	procd_add_reload_trigger "fakehttp"
}

service_stopped() {
	tf_cleanup_nft_table "$SERVICE_NAME"
}
