# luci-app-taoistfuchen/root/etc/init.d/fakehttp
#!/bin/sh /etc/rc.common

# FakeHTTP - disguise outgoing TCP connections as HTTP to evade DPI.
# Managed by luci-app-taoistfuchen.

USE_PROCD=1
START=99
STOP=10

PROG="/usr/bin/fakehttp"
LOG_FILE="/tmp/fakehttp.log"

# Netfilter queue number. MUST differ from FakeSIP (8971) so that both
# services can run at the same time without fighting over the same queue.
QUEUE_NUM=8970

start_service() {
	local enabled interface hostname

	config_load fakehttp
	config_get_bool enabled  main enabled  0
	config_get      interface main interface ''
	config_get      hostname  main hostname  'www.speedtest.net'

	# Disabled in config -> do nothing (this also lets a reload stop it).
	[ "$enabled" = "1" ] || return 0

	if [ -z "$interface" ]; then
		logger -t fakehttp "no network interface configured, not starting"
		return 0
	fi

	[ -n "$hostname" ] || hostname="www.speedtest.net"

	[ -x "$PROG" ] || {
		logger -t fakehttp "binary $PROG missing or not executable"
		return 1
	}

	procd_open_instance fakehttp
	# -i  bind to the chosen interface
	# -h  HTTP obfuscation hostname
	# -s  silent mode: do not log client IP addresses (saves RAM, protects privacy)
	# -n  distinct netfilter queue number (coexist with FakeSIP)
	# -w  write the run log to a file in /tmp (RAM); rotated by a cron job
	# No -d: procd supervises the process in the foreground.
	procd_set_param command "$PROG" -i "$interface" -h "$hostname" -s -n "$QUEUE_NUM" -w "$LOG_FILE"
	procd_set_param respawn
	procd_set_param stderr 1
	procd_close_instance
}

service_triggers() {
	# Reload automatically when the LuCI page saves the config.
	procd_add_reload_trigger "fakehttp"
}
