# luci-app-taoistfuchen/root/etc/init.d/fakesip
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

PROG="/usr/bin/fakesip"
COMMON="/usr/share/taoistfuchen/service-common.sh"

. "$COMMON"

fs_error() {
	logger -t fakesip -- "$*"
}

fs_validate_interface() {
	local iface="$1"

	[ -n "$iface" ] || return 0
	if ! tf_valid_interface "$iface" || ! tf_interface_exists "$iface"; then
		fs_error "invalid or unavailable network interface: $iface"
		FS_VALID=0
		return 0
	fi
	case " $FS_INTERFACES " in *" $iface "*) return 0 ;; esac
	FS_INTERFACES="${FS_INTERFACES:+$FS_INTERFACES }$iface"
	FS_INTERFACE_COUNT=$((FS_INTERFACE_COUNT + 1))
}

fs_validate_config() {
	local interface_mode direction family port_mode ports payload_mode sip_uri repeat ttl
	local hop_estimation dynamic_ttl_pct log_connections queue_num sibling_queue fwmark fwmask

	FS_VALID=1
	FS_INTERFACE_COUNT=0
	FS_INTERFACES=''
	config_get interface_mode main interface_mode 'selected'
	if [ "$interface_mode" != "selected" ]; then
		fs_error "FakeSIP only permits explicitly selected interfaces"
		FS_VALID=0
	fi
	config_list_foreach main interface fs_validate_interface
	if [ "$FS_INTERFACE_COUNT" -eq 0 ]; then
		fs_error "select at least one existing Internet-facing interface"
		FS_VALID=0
	fi

	config_get direction main direction 'outbound'
	case "$direction" in inbound|outbound|both) ;; *) fs_error "invalid direction"; FS_VALID=0 ;; esac
	config_get family main family 'dual'
	case "$family" in
		ipv4|ipv6|dual) ;;
		*) fs_error "invalid address family"; FS_VALID=0 ;;
	esac

	config_get port_mode main port_mode 'exclude'
	config_get ports main ports '53'
	case "$port_mode" in
		include|exclude)
			tf_valid_port_spec "$ports" || { fs_error "invalid UDP port specification"; FS_VALID=0; }
			;;
		all) ;;
		*) fs_error "invalid port filter mode"; FS_VALID=0 ;;
	esac

	config_get payload_mode main payload_mode 'auto'
	config_get sip_uri main sip_uri ''
	case "$payload_mode" in
		auto) ;;
		custom_uri) tf_valid_sip_uri "$sip_uri" || { fs_error "invalid custom SIP URI"; FS_VALID=0; } ;;
		*) fs_error "invalid SIP URI mode"; FS_VALID=0 ;;
	esac

	config_get repeat main repeat '2'
	tf_valid_uint_range "$repeat" 1 10 || { fs_error "repeat must be 1-10"; FS_VALID=0; }
	config_get ttl main ttl '3'
	tf_valid_uint_range "$ttl" 1 255 || { fs_error "TTL must be 1-255"; FS_VALID=0; }
	config_get_bool hop_estimation main hop_estimation 1
	config_get dynamic_ttl_pct main dynamic_ttl_pct '0'
	case "$dynamic_ttl_pct" in
		0) ;;
		*) tf_valid_uint_range "$dynamic_ttl_pct" 1 99 || { fs_error "dynamic TTL percentage must be 0-99"; FS_VALID=0; } ;;
	esac
	if [ "$hop_estimation" != "1" ] && [ "$dynamic_ttl_pct" != "0" ]; then
		fs_error "dynamic TTL requires hop estimation"
		FS_VALID=0
	fi
	config_get queue_num main queue_num '8971'
	tf_valid_uint_range "$queue_num" 1 65535 || { fs_error "NFQUEUE number must be 1-65535"; FS_VALID=0; }
	if [ "$(uci -q get fakehttp.main.enabled)" = "1" ]; then
		sibling_queue="$(uci -q get fakehttp.main.queue_num)"
		if tf_valid_uint_range "$queue_num" 1 65535 &&
		   tf_valid_uint_range "$sibling_queue" 1 65535 &&
		   [ "$queue_num" -eq "$sibling_queue" ]; then
			fs_error "NFQUEUE $queue_num is already assigned to enabled FakeHTTP"
			FS_VALID=0
		fi
	fi
	config_get fwmark main fwmark '0x10000'
	config_get fwmask main fwmask '0x10000'
	tf_valid_u32 "$fwmark" || { fs_error "invalid fwmark"; FS_VALID=0; }
	tf_valid_u32 "$fwmask" || { fs_error "invalid fwmark mask"; FS_VALID=0; }
	if tf_valid_u32 "$fwmark" && tf_valid_u32 "$fwmask"; then
		tf_mark_fits_mask "$fwmark" "$fwmask" || { fs_error "fwmark must fit inside its mask"; FS_VALID=0; }
	fi
	config_get_bool log_connections main log_connections 0
	case "$log_connections:$hop_estimation" in
		0:0|0:1|1:0|1:1) ;;
		*) fs_error "invalid boolean option"; FS_VALID=0 ;;
	esac

	[ "$FS_VALID" = "1" ]
}

fs_append_interface() {
	local iface="$1"

	[ -n "$iface" ] || return 0
	case " $FS_APPENDED_INTERFACES " in *" $iface "*) return 0 ;; esac
	FS_APPENDED_INTERFACES="${FS_APPENDED_INTERFACES:+$FS_APPENDED_INTERFACES }$iface"
	procd_append_param command -i "$iface"
	procd_append_param netdev "$iface"
}

start_service() {
	local enabled direction family port_mode ports payload_mode sip_uri repeat ttl
	local hop_estimation dynamic_ttl_pct log_connections queue_num fwmark fwmask

	tf_cleanup_nft_table fakesip
	config_load fakesip
	config_get_bool enabled main enabled 0
	[ "$enabled" = "1" ] || return 0
	[ -x "$PROG" ] || { fs_error "binary is missing or not executable: $PROG"; return 1; }
	fs_validate_config || return 1

	config_get direction main direction 'outbound'
	config_get family main family 'dual'
	config_get port_mode main port_mode 'exclude'
	config_get ports main ports '53'
	config_get payload_mode main payload_mode 'auto'
	config_get sip_uri main sip_uri ''
	config_get repeat main repeat '2'
	config_get ttl main ttl '3'
	config_get_bool hop_estimation main hop_estimation 1
	config_get dynamic_ttl_pct main dynamic_ttl_pct '0'
	config_get_bool log_connections main log_connections 0
	config_get queue_num main queue_num '8971'
	config_get fwmark main fwmark '0x10000'
	config_get fwmask main fwmask '0x10000'

	procd_open_instance fakesip
	procd_set_param command "$PROG"
	FS_APPENDED_INTERFACES=''
	config_list_foreach main interface fs_append_interface
	# The Droid-MAX baseline used by FakeSIP 0.9.4 checks PACKET_HOST against
	# `outbound` and
	# PACKET_OUTGOING against `inbound`. Compensate for that upstream inversion
	# so the UCI/LuCI direction retains its router-facing meaning.
	case "$direction" in
		inbound) procd_append_param command -1 ;;
		outbound) procd_append_param command -0 ;;
		both) procd_append_param command -0 -1 ;;
	esac
	case "$family" in
		ipv4) procd_append_param command -4 ;;
		ipv6) procd_append_param command -6 ;;
		dual) procd_append_param command -4 -6 ;;
	esac
	case "$port_mode" in
		include) procd_append_param command -p "$ports" ;;
		exclude) procd_append_param command -P "$ports" ;;
	esac
	[ "$payload_mode" = "custom_uri" ] && procd_append_param command -u "$sip_uri"
	procd_append_param command -r "$repeat" -t "$ttl" -n "$queue_num" -m "$fwmark" -x "$fwmask"
	[ "$hop_estimation" = "1" ] || procd_append_param command -g
	[ "$dynamic_ttl_pct" = "0" ] || procd_append_param command -y "$dynamic_ttl_pct"
	[ "$log_connections" = "1" ] || procd_append_param command -s
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param respawn 3600 5 5
	procd_set_param term_timeout 5
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "fakesip"
}

service_stopped() {
	tf_cleanup_nft_table fakesip
}
