# luci-app-taoistfuchen/root/etc/init.d/taoistfuchen-boot-delay
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=98
STOP=09

COMMON="${TF_COMMON:-/usr/share/taoistfuchen/service-common.sh}"
RUNNER="${TF_BOOT_RUNNER:-/usr/share/taoistfuchen/boot-delay-runner.sh}"
INIT_DIR="${TF_INIT_DIR:-/etc/init.d}"

. "$COMMON"

tf_boot_scheduler_error() {
	logger -t taoistfuchen-boot-delay -- "$*" 2>/dev/null || true
}

tf_schedule_boot_service() {
	local service="$1" fallback="$2" enabled delay state phase token deadline now result

	config_load "$service"
	config_get_bool enabled main enabled 0
	config_get delay main boot_delay "$fallback"
	delay="$(tf_boot_delay_value "$delay" "$fallback")" || return 1
	now="$(tf_boot_uptime)" || return 1

	tf_lifecycle_lock_acquire "$service" || return 1
	state="$(tf_boot_state_get "$service" 2>/dev/null || true)"
	if [ "$enabled" != 1 ]; then
		case "$state" in
			wait\ *) tf_boot_state_set "$service" cancel disabled 0 || result=1 ;;
		esac
		tf_lifecycle_lock_release "$service" || return 1
		return ${result:-0}
	fi

	case "$state" in
		manual\ *|done\ *|cancel\ *)
			tf_lifecycle_lock_release "$service" || return 1
			return 0
			;;
		wait\ *)
			set -- $state
			phase="$1"
			token="$2"
			deadline="$3"
			if [ "$deadline" != "$delay" ]; then
				token="$service.$now.$$"
				deadline="$delay"
				tf_boot_state_set "$service" wait "$token" "$deadline" || {
					tf_lifecycle_lock_release "$service" || true
					return 1
				}
			fi
			;;
		*)
			token="$service.$now.$$"
			deadline="$delay"
			tf_boot_state_set "$service" wait "$token" "$deadline" || {
				tf_lifecycle_lock_release "$service" || true
				return 1
			}
			;;
	esac
	tf_lifecycle_lock_release "$service" || return 1

	if [ "$now" -ge "$deadline" ]; then
		"$INIT_DIR/$service" boot_expired "$token" "$deadline" || true
		return 0
	fi

	procd_open_instance "wait-$service"
	procd_set_param command "$RUNNER" "$service" "$token" "$deadline"
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param term_timeout 2
	procd_close_instance
}

start_service() {
	local result=0

	[ -x "$RUNNER" ] || {
		tf_boot_scheduler_error "runner is missing or not executable: $RUNNER"
		return 1
	}
	tf_schedule_boot_service fakesip 40 || result=1
	tf_schedule_boot_service fakehttp 60 || result=1
	return "$result"
}

tf_cancel_boot_waits() {
	local service state result=0

	for service in fakesip fakehttp; do
		tf_lifecycle_lock_acquire "$service" || { result=1; continue; }
		state="$(tf_boot_state_get "$service" 2>/dev/null || true)"
		case "$state" in
			wait\ *) tf_boot_state_set "$service" cancel scheduler-stop 0 || result=1 ;;
		esac
		tf_lifecycle_lock_release "$service" || result=1
	done
	return "$result"
}

stop_service() {
	# rc.common calls this before procd_kill, so no runner can win a wakeup race
	# after an explicit scheduler stop has begun.
	tf_cancel_boot_waits
}

service_stopped() {
	# Idempotent defense for direct procd removal paths.
	tf_cancel_boot_waits
}
