# luci-app-taoistfuchen/root/etc/uci-defaults/99_taoistfuchen
#!/bin/sh

# ---- Custom Logo storage ----
mkdir -p /etc/taoistfuchen
chmod 755 /etc/taoistfuchen

# ---- Make sure the bundled binaries are executable ----
[ -f /usr/bin/fakehttp ] && chmod +x /usr/bin/fakehttp
[ -f /usr/bin/fakesip ]  && chmod +x /usr/bin/fakesip

# ---- Best-effort: pre-fill the WAN egress interface so the user rarely
#      has to figure it out. The default route's interface is the device
#      traffic leaves through. Only set it if not already configured. ----
wan_dev="$(awk '$2 == "00000000" { print $1; exit }' /proc/net/route 2>/dev/null)"
if [ -n "$wan_dev" ]; then
	if [ -z "$(uci -q get fakehttp.main.interface)" ]; then
		uci -q set fakehttp.main.interface="$wan_dev"
	fi
	if [ -z "$(uci -q get fakesip.main.interface)" ]; then
		uci -q set fakesip.main.interface="$wan_dev"
	fi
	uci -q commit fakehttp
	uci -q commit fakesip
fi

# ---- Enable Custom Logo service ----
if [ -f /etc/init.d/taoistfuchen ]; then
	chmod +x /etc/init.d/taoistfuchen
	/etc/init.d/taoistfuchen enable >/dev/null 2>&1 || true
	/etc/init.d/taoistfuchen start  >/dev/null 2>&1 || true
fi

# ---- Register FakeHTTP / FakeSIP with procd. They self-gate on the
#      'enabled' option, so 'start' is a harmless no-op while disabled,
#      but it registers the config-reload triggers immediately. ----
for svc in fakehttp fakesip; do
	if [ -f /etc/init.d/$svc ]; then
		chmod +x /etc/init.d/$svc
		/etc/init.d/$svc enable >/dev/null 2>&1 || true
		/etc/init.d/$svc start  >/dev/null 2>&1 || true
	fi
done

# ---- Install a cron job that trims the FakeHTTP/FakeSIP logs in /tmp every
#      30 minutes so they can never fill RAM. Re-added idempotently. ----
CRONTAB="/etc/crontabs/root"
mkdir -p /etc/crontabs
touch "$CRONTAB"
sed -i '\#taoistfuchen/logclean#d' "$CRONTAB" 2>/dev/null || true
echo "*/30 * * * * /bin/sh /usr/share/taoistfuchen/logclean.sh" >> "$CRONTAB"
/etc/init.d/cron enable  >/dev/null 2>&1 || true
/etc/init.d/cron restart >/dev/null 2>&1 || true

exit 0
