# luci-app-taoistfuchen/root/etc/hotplug.d/iface/99-taoistfuchen
#!/bin/sh

# Keep NFQUEUE bindings in sync with WAN link changes. Hotplug scripts are
# sourced by netifd, so this file deliberately falls through without exit.

case "${ACTION:-}" in
	ifup|ifdown)
		for tf_service in fakehttp fakesip; do
			[ -x "/etc/init.d/$tf_service" ] || continue
			[ "$(uci -q get "$tf_service.main.enabled")" = "1" ] || continue

			tf_matched=0
			tf_mode="$(uci -q get "$tf_service.main.interface_mode")"
			if [ "$tf_service" = "fakehttp" ] && [ "$tf_mode" = "all" ]; then
				tf_matched=1
			elif [ -z "${DEVICE:-}" ]; then
				# Some ifdown events omit DEVICE. Reloading the enabled service is
				# safer than retaining rules bound to a vanished WAN device.
				tf_matched=1
			else
				for tf_iface in $(uci -q get "$tf_service.main.interface"); do
					[ "$tf_iface" = "$DEVICE" ] && tf_matched=1
				done
			fi

			if [ "$tf_matched" = "1" ]; then
				/etc/init.d/$tf_service reload >/dev/null 2>&1 || true
			fi
		done
		unset tf_service tf_matched tf_mode tf_iface
		;;
esac
