# luci-app-taoistfuchen/root/etc/hotplug.d/iface/99-taoistfuchen
#!/bin/sh
# Boot-time hardening for FakeHTTP / FakeSIP auto-start.
#
# When the device a service binds to comes up, make sure the service is
# running. This covers slow WAN links (e.g. PPPoE) where the interface may
# not exist yet at START=99 boot time.
#
# IMPORTANT: we only START the service when it is NOT already running. We do
# NOT restart a healthy instance, because some links (PPPoE renegotiation,
# IPv6 prefix refreshes, etc.) emit ifup events repeatedly; restarting on each
# of them would put the service into a stop/start loop. A running instance
# keeps working across those events on its own (its nft rules match by
# interface name), and procd re-spawns it if it ever exits.
#
# hotplug.d scripts are sourced, so this script avoids calling 'exit' (which
# would abort the whole hotplug run). It ends naturally instead.

if [ "$ACTION" = "ifup" ]; then
	for svc in fakehttp fakesip; do
		[ -x "/etc/init.d/$svc" ] || continue

		# Only services enabled in their config.
		[ "$(uci -q get "$svc.main.enabled")" = "1" ] || continue

		# Only when the just-upped device is the one the service binds to.
		iface="$(uci -q get "$svc.main.interface")"
		[ -n "$iface" ] || continue
		[ "$DEVICE" = "$iface" ] || continue

		# Start ONLY if it is not already running (no restart loops).
		if ! pidof "$svc" >/dev/null 2>&1; then
			/etc/init.d/$svc start >/dev/null 2>&1 &
		fi
	done
fi
